# etc/s6-overlay/s6-rc.d/configure-ipsec/up
/etc/s6-overlay/scripts/configure-ipsec.sh
