# etc/s6-overlay/s6-rc.d/configure-webserver/up
ln -sfn /etc/swanctl/client.mobileconfig /var/www/client.mobileconfig
