# etc/s6-overlay/s6-rc.d/strongswan/run
#!/command/execlineb -P
ipsec start --nofork
